# supervisord.conf
[supervisord]
nodaemon=true

[program:squid]
command=/bin/bash /start-squid.sh

